# roles/cpu_temp_mon_bare/templates/prometheus.service.j2
[Unit]
Description=Prometheus Monitoring
Wants=network-online.target
After=network-online.target

[Service]
ExecStart={{ prometheus_service_exec_start }}
Restart=always

[Install]
WantedBy=multi-user.target
